# tools/backlight-tooler.service
[Unit]
Description=Automated Backlight Setting

[Service]
Type=oneshot
ExecStart=/usr/bin/backlight-tooler
