# tools/BacklightTooler.service
[Unit]
Description=Automated Backlight Setting

[Service]
Type=oneshot
ExecStart=/usr/bin/BacklightTooler auto 12
