# tools/backlight-tooler.timer
[Unit]
Description=Timer to Control the Backlight Automation

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
